# etc/cron.hourly/dshield
#!/bin/bash

####
#
#  hourly checking with dshield. to be expanded for update check
#
####

userid=`id -u`
if [ ! "$userid" = "0" ]; then
   echo "you have to run this script as root. eg."
   echo "  sudo dshield (or as root from cron)"
   exit
fi

# check if automatic updates are allowed
if [ ! -f /srv/dshield/auto-update-ok ] ; then
   exit
fi

version=0.0
source <(grep = /etc/dshield.ini | sed 's/ *= */=/g')
nonce=`openssl rand -hex 10`
hash=`echo -n $email:$apikey | openssl dgst -hmac $nonce -sha512 -hex | cut -f2 -d'=' | tr -d ' '`
# TODO: urlencode($user)
user=`echo $email | sed 's/+/%2b/' | sed 's/@/%40/'`
newversion=`curl -s https://isc.sans.edu/api/checkapikey/$user/$nonce/$hash | grep '<result>ok</result>' | grep '\<version\>' | sed 's/.*<version>//' | sed 's/<\/version>.*//'`
majornewversion=`echo $newversion | cut -f1 -d'.'`
majoroldversion=`echo $version | cut -f1 -d'.'`
minornewversion=`echo $newversion | cut -f2 -d'.'`
minoroldversion=`echo $version | cut -f2 -d'.'`
if [ ! -z $PS1 ]; then
    echo "Current Version: $version ; New Version: $newversion"
fi


if [ "$majoroldversion" -lt "$majornewversion" ]; then
    if [ ! -z $PS1 ]; then
       echo "UPDATE"
    fi
    git -C $progdir pull -q
    $progdir/update.sh
    exit
fi
if [ $majoroldversion -eq $majornewversion ]; then
    if [ $minoroldversion -lt $minornewversion ]; then
	if [ ! -z $PS1 ]; then
		echo "UPDATE"
	fi
	
	git -C $progdir  pull -q
	$progdir/update.sh
	exit
    fi
fi
if [ -z $PS1 ]; then
    exit
else
    echo "NO UPDATE"
fi
