# etc/network/if-pre-up.d/dshield
#!/bin/sh
iptables-restore /etc/network/iptables
if [ -f /etc/network/iptables.local ]; then
    iptables-restore -n /etc/network/iptables.local
fi
exit 0
