# etc/network/if-pre-up.d/dshield
#!/bin/sh
iptables-restore < /etc/network/iptables
exit 0